# ops/systemd/wiki-engine-runner.timer
[Unit]
Description=Run Wiki Engine Codex task runner every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Persistent=true
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
